# etc/profile.d/clearos.sh
# Add new path to support version wrappers (e.g. PHP Engines)
export PATH="/usr/clearos/bin:${PATH}"
